# arch-system/base-system.install
post_upgrade() {
	# in case it wasn't already run by glibc
	locale-gen

	systemctl enable --now avahi-daemon.service
}

post_install() {
	post_upgrade
}

pre_remove() {
	systemctl disable avahi-daemon.service
}
